# helm/cluster-aws/files/etc/systemd/system/cloud-config.service.d/http-proxy.conf
[Service]
ExecStart=/bin/bash /opt/cloud-config.sh
